# deploy/supervisor/facturasapp.conf
[program:facturasapp]
command= /home/jair/facturasapp/deploy/sh/server.sh
user=jair
stdout_logfile= /home/jair/facturasapp/logs/supervisor.log
redirect_stderr= true
autostart= true
autorestart= true
environment=LANG= en_US.UTF-8,LC_ALL=en_US.UTF-8
